# ticket-flags.service
[Unit]
Description=Traitement des flags

[Service]
Type=oneshot
ExecStart=/opt/traitement-presences/flags-loop.sh
