# ticket-flags.timer
[Unit]
Description=Run flags service every minute

[Timer]
OnCalendar=*:0/01

[Install]
WantedBy=timers.target