# utils/ticket-flags.service
[Unit]
Description=Traitement des flags

[Service]
Type=oneshot
ExecStart=/opt/traitement-presences/utils/flags-loop.sh
